# hyperos/airootfs/etc/systemd/system/hyperos-gamed.service
[Unit]
Description=Hyper OS gaming profile daemon
After=multi-user.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hyperos-gamed
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
