# hyperos/airootfs/etc/profile.d/hyperos-prompt.sh
export PS1='\[\e[38;2;110;193;255m\]\u@hyper\[\e[0m\]:\[\e[38;2;139;92;246m\]\w\[\e[0m\]\$ '
